# etc/pdudaemon.service
[Unit]
Description=TCP Listening daemon to accept PDU requests
After=postgresql.service

[Service]
ExecStart=/usr/sbin/pdudaemon-start -j
Type=forking
PIDFile=/var/run/pdudaemon.pid

[Install]
WantedBy=network.target
